# roles/certificates/templates/cert-renew.timer.j2
## {{ ansible_managed }}
# SPDX-License-Identifier: MIT-0
#
# Systemd timer for automatic certificate renewal
# Runs the cert-renew.service twice daily

[Unit]
Description=Certificate Renewal Timer
Documentation=https://github.com/yourusername/jlira.web_server
Requires=cert-renew.service

[Timer]
# Run twice daily (at 00:00 and 12:00)
OnCalendar=daily
OnCalendar=12:00

# Add randomized delay to avoid thundering herd
RandomizedDelaySec=3600

# Run on boot if the last run was missed
Persistent=true

[Install]
WantedBy=timers.target
